# base/services.d/xvfb/run
#!/usr/bin/with-contenv bash

Xvfb "$DISPLAY" -screen 0 "$DISPLAY_RESOLUTION"
