# base/services.d/x11vnc/run
#!/usr/bin/with-contenv bash

x11vnc -display "$DISPLAY" -xkb -forever -passwd "$VNC_PASSWORD" -rfbport "$VNC_PORT"
