# openbox/services.d/openbox-session/run
#!/usr/bin/with-contenv bash

openbox-session
