# base/archlinux/services.d/vncserver/run
#!/usr/bin/with-contenv bash

vncserver
